# supervisord.conf
[supervisord]
nodaemon=true
user=kanshi
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:ustreamer]
command=/bin/bash -c 'source /app/config/ustreamer.conf && exec /usr/local/bin/ustreamer $USTREAMER_OPTS'
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:capture_uploader]
command=python3 /app/src/capture_uploader.py
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:web_server]
command=python3 /app/src/web_server.py
autostart=true
autorestart=true
startretries=5
startsecs=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
